# xfce-core-meta/files/etc/profile.d/xfce-core-meta-qt5ct.sh
[ ! -z "$XDG_SESSION_DESKTOP" ] && [ "$XDG_SESSION_DESKTOP" == "xfce" ] && export QT_QPA_PLATFORMTHEME=qt5ct
